# auto_down.cron
*/1 * * * *  /root/OneGuyAPI/auto_down.sh